# tests/cron/fs/etc/cron.d/backups
* * * * * cron tar -zcf /var/backups/home.tgz /home/
