# experiments/ubuntu/cgc/xinetd.d/cgc-model
service cadet_00001 
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cadet_00001.sh
        socket_type = stream
        protocol = tcp
        user = opam 
        port = 8001 
        flags = REUSE
        per_source = UNLIMITED	
}
