# experiments/ubuntu/cgc/xinetd.d/cgc-model
service cadet_00001
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cadet_00001.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8001
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00007
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00007.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8002
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00046
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00046.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8003
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00047
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00047.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8004
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00048
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00048.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8005
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00051
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00051.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8006
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00054
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00054.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8007
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00055
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00055.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8008
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00057
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00057.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8009
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00058
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00058.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8010
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00061
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00061.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8011
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00064
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00064.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8012
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00065
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00065.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8013
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00066
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00066.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8014
        flags = REUSE
        per_source = UNLIMITED
}

service cromu_00070
{
        disable = no
        type = UNLISTED
        wait = no
        server = /uIDS/experiments/ubuntu/cgc/scripts/cromu_00070.sh
        socket_type = stream
        protocol = tcp
        user = opam
        port = 8015
        flags = REUSE
        per_source = UNLIMITED
}
